# ci/update-nc.service
[Unit]
Description=NC CI update.

[Service]
Type=oneshot
ExecStart=/wagglerw/message-pipeline/ci/update-nc.sh
