# ci/update-nc.timer
[Unit]
Description=NC CI update.

[Timer]
OnBootSec=5min
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
